# base/conf/yaffas-postfix/debian/postinst
#!/bin/sh
# postinst script for bbconf-postfix
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)

		CONF=/etc/postfix
		mkdir -p $CONF

		touch $CONF/ldap-aliases.cf
		touch $CONF/ldap-aliases.cf.db
		touch $CONF/ldap-users.cf
		touch $CONF/ldap-users.cf.db
		touch $CONF/smtp_auth.cf
        postmap $CONF/smtp_auth.cf
        touch $CONF/virtual_users_global
        postmap $CONF/virtual_users_global

		chmod 600 $CONF/smtp_auth.cf
		chmod 600 $CONF/smtp_auth.cf.db

        H=$(/bin/hostname -f)
        sed -e 's/HOSTNAME/'$H'/' -i $CONF/main.cf

        /usr/bin/newaliases

        invoke-rc.d postfix restart
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


