# base/webmin/yaffas-core/debian/postinst
#!/bin/bash
# postinst script for bbwebmin-core
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package
#
# quoting from the policy:
#     Any necessary prompting should almost always be confined to the
#     post-installation script, and should be protected with a conditional
#     so that unnecessary prompting doesn't happen if a package's
#     installation fails and the `postinst' is called with `abort-upgrade',
#     `abort-remove' or `abort-deconfigure'.

case "$1" in
    configure)

	if [[ ! -L /etc/init.d/yaffas ]]; then 
		ln -sf /opt/yaffas/etc/init.d/yaffas.debian /etc/init.d/yaffas
		update-rc.d yaffas defaults
	fi
	
	# generate passwort file for admin logging if no one exists
	if [ ! -f /opt/yaffas/etc/webmin/miniserv.users ]; then
		echo "admin:129TZgbE1H546:0" > /opt/yaffas/etc/webmin/miniserv.users
	fi

	# generate acl file for webmin if no one exists
	if [ ! -f /opt/yaffas/etc/webmin/webmin.acl-global ]; then
        echo "admin: " > /opt/yaffas/etc/webmin/webmin.acl-global
        echo "admin: setup" > /opt/yaffas/etc/webmin/webmin.acl-setup
        ln -sf /opt/yaffas/etc/webmin/webmin.acl-setup /opt/yaffas/etc/webmin/webmin.acl
	fi

	touch /opt/yaffas/etc/webmin/hidden_modules

	# correct permissions
	chmod 600 /opt/yaffas/etc/webmin/miniserv.*

	# restart webmin (stop script can crash, so lets kill it!)
	PID=`ps ax | grep -v grep | grep "/opt/yaffas/webmin/miniserv.pl" | awk {' print $1 '}`	
	if [ "$PID" ]; then
	        kill -9 $PID
	        rm -f /opt/yaffas/var/miniserv.pid
	fi
	/opt/yaffas/etc/webmin/start
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


