# base/webmin/modules/yaffas-module-snmpconf/debian/postinst
#!/bin/bash
# postinst script for bbusermin-bbsetvacation
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package
#

source /opt/yaffas/lib/bbinstall-lib.sh


case "$1" in
    configure)

	MODULE="snmpconf"
	add_webmin_acl $MODULE
	add_license $MODULE ""
	
	if cat /etc/snmp/snmpd.conf| grep com2sec | grep paranoid >/dev/null; then
	    /usr/sbin/update-rc.d -f snmpd remove
	fi

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0



